# PROJECTS/beginner/systemd-persistence-scanner/testdata/systemd/watch-beacon.path
[Unit]
Description=Watch for beacon removal

[Path]
PathExists=/tmp/.beacon
Unit=restore-beacon.service

[Install]
WantedBy=multi-user.target
